# infra/systemd/mitmproxy.service
[Unit]
Description=mitmproxy for Blaulicht per-app VPN
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mitmproxy
Group=mitmproxy
WorkingDirectory=/opt/blaulicht/mitmproxy
# TIP: mitmdump is the headless binary. Adjust paths/ports as needed.
StateDirectory=mitmproxy
# Explicitly set mitmproxy confdir to a systemd-managed writable location
ExecStart=/usr/bin/mitmdump --set confdir=/var/lib/mitmproxy --mode socks5 --listen-host 0.0.0.0 --listen-port 8281 -s /opt/blaulicht/mitmproxy/garage_trigger.py
Restart=on-failure
RestartSec=2s

# Hardening ideas (uncomment and test):
# NoNewPrivileges=true
# PrivateTmp=true
# ProtectSystem=full
# ProtectHome=true

[Install]
WantedBy=multi-user.target
